# config/etcd.ini
[program:etcd]
directory=/usr/local/bin/
command=/usr/local/bin/etcd --config-file=/etc/etcd.conf
autostart=true
autorestart=true
stderr_logfile=/tmp/etcd_stderr.log
stdout_logfile=/tmp/etcd_stdout.log
user = root

