# config/nats.ini
[program:nats]
directory=/usr/local/bin/
command=/usr/local/bin/nats-server -a 192.168.87.237
autostart=true
autorestart=true
stderr_logfile=/tmp/nats_stderr.log
stdout_logfile=/tmp/nats_stdout.log
user = root

